# ops/systemd/exam-pg-backup.service
[Unit]
Description=Exam Platform — daily PostgreSQL backup
After=network-online.target postgresql.service
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=/srv/exam-platform/.env
ExecStart=/srv/exam-platform/ops/backup/pg-backup.sh
ExecStartPost=/srv/exam-platform/ops/backup/uploads-backup.sh
User=root
Nice=10
IOSchedulingClass=idle
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ReadWritePaths=/var/backups/postgres /var/log /tmp

[Install]
WantedBy=multi-user.target
